# centinela_work/olt_monitor/systemd/olt-monitor-collect.service
[Unit]
Description=OLT Optical Monitor — collect snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/olt_monitor
ExecStart=/usr/bin/python3 /opt/olt_monitor/monitor.py run
StandardOutput=append:/var/log/olt-monitor.log
StandardError=append:/var/log/olt-monitor.log
User=root
Nice=10
TimeoutStartSec=10min
